# systemd/owui.service
[Unit]
Description=OpenWebUI Suite (docker compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root
Environment=COMPOSE_DOCKER_CLI_BUILD=1
Environment=DOCKER_BUILDKIT=1
Environment=COMPOSE_PROFILES=extras,ui
ExecStartPre=/bin/bash -lc '/usr/bin/docker network inspect ${OWUI_NETWORK:-root_owui} >/dev/null 2>&1 || /usr/bin/docker network create ${OWUI_NETWORK:-root_owui}'
ExecStart=/usr/bin/docker compose -f /root/compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f /root/compose.prod.yml stop
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
